# meta-mchp/meta-mchp-common/recipes-phosphor/dummy-power/files/dummy-power@.service
[Unit]
Description=Dummy OpenBMC Power Control D-Bus Service

[Service]
ExecStart=/usr/bin/dummy_power
Type=simple

[Install]
WantedBy=multi-user.target
